# src/lib/systemd/system/angusnet-route-fixd.service
[Unit]
Description=AMPR route fixer

[Service]
ExecStart=/usr/local/sbin/angusnet-route-fixd
Restart=always

[Install]
WantedBy=network.target
