# src/lib/systemd/system/angusnet-info-bbs.service
[Unit]
Description=AMPR route fixer

[Service]
ExecStart=/usr/local/sbin/angusnet-info-bbs
Restart=always

[Install]
WantedBy=network.target
