# src/etc/systemd/system/ampr-route-fixd.service
[Unit]
Description=AMPR route fixer

[Service]
ExecStart=/usr/local/sbin/ampr-route-fixd
Restart=always

[Install]
WantedBy=network.target
